# simple.conf
[program:myserver]
command=gunicorn run:app -w 4 -b 0.0.0.0:5000
directory=/home/ubuntu/trueFitSkinCare/front_end/

[supervisord]
logfile=/home/ubuntu/supervisord.log
loglevel=debug
user=root

